# prod/supervisor.d/django-uwsgi.ini
[program:uwsgi]
command=/usr/sbin/uwsgi --ini /app/prod/uwsgi.ini
autostart=true
autorestart=true
startsecs=10
user=www-data
stopsignal=QUIT
stdout_logfile=/app/logs/uwsgi.log
stderr_logfile=/app/logs/uwsgi-error.log
