# prod/supervisor.d/nginx.ini
[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/app/logs/nginx.log
stderr_logfile=/app/logs/nginx-error.log
